# etc/systemd/system/android-mount.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/lib/droidian/mount-android.sh
