# etc/systemd/system/android-mount.service.d/override.conf
[Service]
Environment=ANDROID_DATA_PART=/dev/disk/by-partlabel/userdata
Environment=ANDROID_SYSTEM_PART=/dev/disk/by-partlabel/system
Environment=ANDROID_VENDOR_PART=/dev/disk/by-partlabel/vendor
