# etc/systemd/system/plymouth-start.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/plymouthd --mode=boot --pid-file=/run/plymouth/pid --attach-to-session --renderer-type=frame-buffer
