# usr/lib/systemd/system/backlight-restore.path
[Unit]
Description=Restore backlight brightness on display resume

[Path]
PathChanged=/sys/class/drm/card0-DSI-1/dpms
Unit=backlight-restore.service

[Install]
WantedBy=multi-user.target
