# usr/lib/systemd/system/ppd-gov-handler.path
[Unit]
Description=Monitor power profile changes

[Path]
PathModified=/var/lib/power-profiles-daemon/state.ini

[Install]
WantedBy=multi-user.target
