# etc/systemd/user/pulseaudio.service.d/override.conf
[Unit]
After=lxc@android.service
Wants=lxc@android.service

[Service]
ExecStartPre=
ExecStartPre=/bin/sh -c '/usr/bin/sleep 3'
Restart=on-failure
RestartSec=5
