# etc/systemd/system/power-profiles-daemon.service.d/fake-driver.conf
[Service]
Environment=POWER_PROFILE_DAEMON_FAKE_DRIVER=1
